# redwoodhq/supervisord.conf
[supervisord]
nodaemon=true
user=root
directory=/data/supervisor
childlogdir=/data/supervisor
logfile=supervisord.log

[program:redwoodhq]
command=/rootstart.sh


